# 009/server.service
[Unit]
Description=UDP Server Service

[Service]
ExecStart=/app/server
Restart=always
RestartSec=30
TimeoutStopSec=90
