# 012/server.service
[Unit]
Description=UDP Server Service

[Service]
ExecStart=/app/server ens3
Restart=always
RestartSec=30
TimeoutStopSec=90
